# packages/cortex-ai-github/cortex-ai-github.service
[Unit]
Description=Cortex AI GitHub App
After=network.target

[Service]
Type=simple
User=jamiecraik
WorkingDirectory=/Users/jamiecraik/.Cortex-OS/packages/cortex-ai-github
Environment=NODE_ENV=production
EnvironmentFile=/Users/jamiecraik/.Cortex-OS/packages/cortex-ai-github/.env
ExecStart=/usr/local/bin/npx tsx src/server/start.ts
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cortex-ai-github

[Install]
WantedBy=multi-user.target
