# scripts/hf/migrate.models.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cortexos.migrate-hf-models</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jamiecraik/.Cortex-OS/scripts/hf/migrate_models_to_hub.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>86400</integer><!-- run once per day -->
  <key>StandardOutPath</key>
  <string>/tmp/migrate-hf-models.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/migrate-hf-models.err.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
