# infra/mcp/com.cortexos.mcp.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cortexos.mcp.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jamiecraik/.Cortex-OS/scripts/start-mcp-server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jamiecraik/.Cortex-OS/packages/cortex-mcp</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jamiecraik/Library/Logs/com.cortexos.mcp.server.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jamiecraik/Library/Logs/com.cortexos.mcp.server.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jamiecraik/.local/share/mise/installs/python/3.12.6/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/jamiecraik/.Cortex-OS/packages/cortex-mcp</string>
        <key>UVICORN_LOG_LEVEL</key>
        <string>info</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>HOME</key>
        <string>/Users/jamiecraik</string>
    </dict>
</dict>
</plist>
