# infra/mcp/mcp-server.service
[Unit]
Description=MCP Server for brAInwav Cortex-OS
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=jamiecraik
WorkingDirectory=/Users/jamiecraik/.Cortex-OS/packages/cortex-mcp
Environment=PATH=/Users/jamiecraik/.local/share/mise/installs/python/3.12.6/bin
Environment=PYTHONPATH=/Users/jamiecraik/.Cortex-OS/packages/cortex-mcp
ExecStart=/Users/jamiecraik/.Cortex-OS/scripts/start-mcp-server.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
